# config/sensevoice-vibe.service.example
[Unit]
Description=SenseVoice Resident VAD Daemon
After=default.target
StartLimitIntervalSec=30
StartLimitBurst=20

[Service]
Type=simple
WorkingDirectory=/home/dell/mosim_workspace/work/sensevoice-local-small
EnvironmentFile=-%h/.config/sensevoice-vibe/llm.env
Environment=SENSEVOICE_RESIDENT=1
Environment=SENSEVOICE_STREAM_ACTIVE_ON_START=1
Environment=SENSEVOICE_MODEL=%h/.cache/modelscope/hub/models/iic/SenseVoiceSmall
Environment=SENSEVOICE_LANGUAGE=auto
Environment=SENSEVOICE_FILTER_FILLERS=1
Environment=SENSEVOICE_SPK_ENABLE=1
Environment=SENSEVOICE_SPK_ENROLL_WAV=%h/.local/state/sensevoice-vibe/speaker_enroll.d
Environment=SENSEVOICE_SPK_THRESHOLD=0.52
Environment=SENSEVOICE_SPK_MIN_MS=900
Environment=SENSEVOICE_SPK_MODEL=speechbrain/spkrec-ecapa-voxceleb
Environment=SENSEVOICE_SPK_CACHE_DIR=%h/.cache/sensevoice-vibe/spkrec
Environment=SENSEVOICE_SPK_AGG=topk_mean
Environment=SENSEVOICE_SPK_TOPK=1
Environment=SENSEVOICE_SPK_AUTO_ENROLL=0
Environment=SENSEVOICE_SPK_AUTO_ENROLL_DIR=%h/.local/state/sensevoice-vibe/speaker_enroll.d
Environment=SENSEVOICE_SPK_AUTO_ENROLL_MIN_SCORE=0.70
Environment=SENSEVOICE_SPK_AUTO_ENROLL_MIN_MS=2200
Environment=SENSEVOICE_SPK_AUTO_ENROLL_COOLDOWN_SEC=90
Environment=SENSEVOICE_SPK_AUTO_ENROLL_MAX_TEMPLATES=20
Environment=SENSEVOICE_SPK_ADAPTIVE=1
Environment=SENSEVOICE_SPK_ADAPTIVE_WINDOW=80
Environment=SENSEVOICE_SPK_ADAPTIVE_MIN_SAMPLES=10
Environment=SENSEVOICE_SPK_ADAPTIVE_FLOOR=0.52
Environment=SENSEVOICE_SPK_ADAPTIVE_MARGIN=0.04
Environment=SENSEVOICE_SPK_PRUNE_OUTLIERS=1
Environment=SENSEVOICE_SPK_PRUNE_KEEP=10
Environment=SENSEVOICE_PARTIAL_STRATEGY=stable2
Environment=SENSEVOICE_EMIT_PARTIAL=0
Environment=SENSEVOICE_AUTO_ENTER=0
Environment=SENSEVOICE_PREFER_WTYPE=0
Environment=SENSEVOICE_PASTE_KEY=ctrl_v
Environment=SENSEVOICE_YDOTOOL_KEY_DELAY_MS=35
Environment=SENSEVOICE_CLIPBOARD_SETTLE_SEC=0.08
Environment=SENSEVOICE_PASTE_PRE_DELAY_SEC=0.08
Environment=SENSEVOICE_CLIPBOARD_RESTORE=1
Environment=SENSEVOICE_CLIPBOARD_RESTORE_DELAY_SEC=0.12
Environment=SENSEVOICE_CLIPBOARD_VERIFY=1
Environment=SENSEVOICE_INJECT_ACK_TIMEOUT_SEC=1.2
Environment=SENSEVOICE_DEBUG_INJECT=1
Environment=SENSEVOICE_CLEAR_BEFORE_REPLACE=0
Environment=SENSEVOICE_STREAM_INDICATOR=none
Environment=SENSEVOICE_STREAM_VAD_AGGRESSIVENESS=3
Environment=SENSEVOICE_STREAM_START_MS=240
Environment=SENSEVOICE_STREAM_FRAME_MS=20
Environment=SENSEVOICE_STREAM_PRE_ROLL_MS=700
Environment=SENSEVOICE_STREAM_ENDPOINT_MS=900
Environment=SENSEVOICE_STREAM_MAX_SEGMENT_MS=30000
Environment=SENSEVOICE_STREAM_MIN_SEGMENT_MS=850
Environment=SENSEVOICE_STREAM_PARTIAL_INTERVAL_MS=280
Environment=SENSEVOICE_STREAM_MIN_PARTIAL_MS=1300
Environment=SENSEVOICE_RETENTION_KEEP_RECENT=20
Environment=SENSEVOICE_LOG_KEEP_LINES=120
Environment=SENSEVOICE_POST_LLM_ENABLE=1
Environment=SENSEVOICE_POST_LLM_MODEL=DeepSeek-V3.1-Terminus
Environment=SENSEVOICE_POST_LLM_FALLBACK_MODEL=DeepSeek-V3.1-Terminus
Environment=SENSEVOICE_POST_LLM_MODE=polish_coding_aggressive
Environment=SENSEVOICE_POST_LLM_TIMEOUT_MS=1800
Environment=SENSEVOICE_POST_LLM_MAX_TOKENS=72
Environment=SENSEVOICE_POST_LLM_TEMPERATURE=0
Environment=SENSEVOICE_POST_LLM_STRICT=0
Environment=SENSEVOICE_POST_LLM_CIRCUIT_MAX_FAILS=4
Environment=SENSEVOICE_POST_LLM_CIRCUIT_COOLDOWN_SEC=25
Environment=SENSEVOICE_POST_LLM_HARD_COOLDOWN_SEC=300
Environment=SENSEVOICE_POST_LLM_RETRY_ON_TIMEOUT=1
Environment=SENSEVOICE_POST_LLM_RETRY_BACKOFF_MS=80
Environment=SENSEVOICE_POST_LLM_MODEL_AUTO=1
Environment=SENSEVOICE_POST_LLM_MODEL_PROBE_TIMEOUT_MS=450
Environment=SENSEVOICE_POST_LLM_MIN_CHARS=5
Environment=SENSEVOICE_POST_LLM_CACHE_TTL_SEC=300
Environment=SENSEVOICE_POST_LLM_CACHE_MAX_ENTRIES=120
Environment=SENSEVOICE_POST_LLM_DYNAMIC_MAX_TOKENS=1
Environment=SENSEVOICE_POST_LLM_OUTPUT_TOKEN_FACTOR=0.7
Environment=SENSEVOICE_PROJECT_LEXICON_ENABLE=1
Environment=SENSEVOICE_PROJECT_ROOT=%h/mosim_workspace
Environment=SENSEVOICE_PROJECT_LEXICON_MAX_TERMS=2500
Environment=SENSEVOICE_PROJECT_LEXICON_HINT_LIMIT=16
Environment=SENSEVOICE_PROJECT_LEXICON_MIN_TERM_LEN=3
Environment=SENSEVOICE_PROJECT_LEXICON_EXTRA_FILE=/home/dell/mosim_workspace/work/sensevoice-local-small/hotwords_coding_zh.txt
Environment=SENSEVOICE_CONF_ROUTE_ENABLE=1
Environment=SENSEVOICE_CONF_ROUTE_HIGH=0.42
Environment=SENSEVOICE_CONF_ROUTE_LOW=0.30
Environment=SENSEVOICE_LEARN_ENABLE=1
Environment=SENSEVOICE_LEARN_STORE=%h/.local/state/sensevoice-vibe/correction_memory.json
Environment=SENSEVOICE_LEARN_MIN_HITS=2
Environment=SENSEVOICE_LEARN_MAX_RULES=320
Environment=SENSEVOICE_COMPARE_LOG_ENABLE=1
Environment=SENSEVOICE_COMPARE_LOG_FILE=%h/.local/state/sensevoice-vibe/post_compare.jsonl
Environment=SENSEVOICE_COMPARE_LOG_KEEP_LINES=300
Environment=SENSEVOICE_STATE_LOG=%h/.local/state/sensevoice-vibe/stream_vad.log
ExecStartPre=/bin/bash /home/dell/mosim_workspace/work/sensevoice-local-small/retention_housekeep.sh
ExecStart=/home/dell/mosim_workspace/work/sensevoice-local-small/.venv/bin/python /home/dell/mosim_workspace/work/sensevoice-local-small/stream_vad_realtime.py --resident --language auto --indicator none
Restart=always
RestartSec=0.8
TimeoutStopSec=3

[Install]
WantedBy=default.target
