# config/sensevoice-vibe.service.example
[Unit]
Description=SenseVoice Resident VAD Daemon
After=default.target
StartLimitIntervalSec=30
StartLimitBurst=20

[Service]
Type=simple
# __INSTALL_DIR__ 由 install_autostart_service.sh 自动替换为实际安装路径
WorkingDirectory=__INSTALL_DIR__
# 所有配置从 llm.env 读取（单一配置源）
EnvironmentFile=-%h/.config/sensevoice-vibe/llm.env
Environment=PYTHONPATH=__INSTALL_DIR__
Environment=SENSEVOICE_RESIDENT=1
Environment=SENSEVOICE_STREAM_ACTIVE_ON_START=1
Environment=SENSEVOICE_STATE_LOG=%h/.local/state/sensevoice-vibe/stream_vad.log
Environment=SENSEVOICE_COMPARE_LOG_FILE=%h/.local/state/sensevoice-vibe/post_compare.jsonl
Environment=SENSEVOICE_LEARN_STORE=%h/.local/state/sensevoice-vibe/correction_memory.json
ExecStartPre=/bin/bash __INSTALL_DIR__/retention_housekeep.sh
ExecStart=__INSTALL_DIR__/.venv/bin/python __INSTALL_DIR__/stream_vad_realtime.py --resident --language auto --indicator none
Restart=always
RestartSec=0.8
TimeoutStopSec=3

[Install]
WantedBy=default.target
